# services/wowos-install-desktop-once.service
[Unit]
Description=wowOS install desktop once on first boot
After=network-online.target
Wants=network-online.target
ConditionPathExists=!/var/lib/wowos/.desktop_installed

[Service]
Type=oneshot
ExecStart=/opt/wowos/scripts/install_desktop_once.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
