# services/wowos-install-desktop-once.service
[Unit]
# 方案2: 首启安装桌面，标记文件 .desktop_installed，脚本 /opt/wowos/scripts/install_desktop_once.sh
Description=wowOS install desktop once on first boot
After=network.target
ConditionPathExists=!/var/lib/wowos/.desktop_installed

[Service]
Type=oneshot
ExecStart=/opt/wowos/scripts/install_desktop_once.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
